# systemd/file-backup.timer.template
[Unit]
Description=Backup Suite file backup timer

[Timer]
OnCalendar=__FILE_BACKUP_ONCALENDAR__
Persistent=true
Unit=file-backup.service

[Install]
WantedBy=timers.target
